# units/wsl-theme-sync.service
[Unit]
Description=Mirror Windows light/dark theme into GNOME
Documentation=https://github.com/techneut92/wsl-gnome-rdp-installer
# After=default.target so XDG_CURRENT_DESKTOP + the user dbus session
# are settled — gsettings doesn't fault if dconf isn't reachable, but
# the value won't get written cleanly.
After=default.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/wsl-theme-sync
# Strict ceiling — the script wraps reg.exe + each gsettings call in
# `timeout 3`, so a healthy run completes in well under a second. 15s
# is "something is genuinely wedged, kill rather than queue the next
# .timer fire behind a hung instance." Without this, a hung run
# leaves the unit in `activating start` indefinitely, blocking
# subsequent fires and eventually starving dbus traffic on the user
# bus.
TimeoutStartSec=15

[Install]
WantedBy=default.target
